# examples/services/rusnapshot-monthly.timer
[Unit]
Description=Take monthly snapshots with rusnapshot

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target